# roles/srvbak/files/0am.cron
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

5 0 * * * root  cd / && run-parts --report /etc/cron.0am
