# roles/srvbak/files/srvbak.cron
#!/bin/sh

set -e
/opt/src/srvbak/srvbak.bash /opt/src/srvbak/srvbakrc 2>&1

# --
